# packaging/systemd/auspexai-worker.service
; AuspexAI worker — systemd user service unit template.
;
; Install: `systemctl --user enable --now auspexai-worker.service`
;
; M2: heartbeat loop. M3+ adds assignment pull. The hardening directives
; below are PHASE-1 PERMISSIVE per §5.17 — Phase 2 tightens defaults via
; config without rewriting the unit.

[Unit]
Description=AuspexAI volunteer worker (Phase 1)
Documentation=https://github.com/auspexai/worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/auspexai-worker daemon
Restart=on-failure
RestartSec=10
; Hardening per §5.17 — applied at the daemon (trusted) tier. Phase 1
; permissive defaults; Phase 2 narrows further via config without altering
; the unit body.
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=strict
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
; §41(a): the daemon puts per-unit cgroup-v2 caps (memory.max/pids.max) around
; the runner, so the control-group tree must stay WRITABLE (=no) and delegated
; (Delegate=yes below). ProtectControlGroups=yes would remount it read-only.
ProtectControlGroups=no
RestrictRealtime=true
; @mount is REQUIRED on top of @system-service: the daemon spawns bubblewrap,
; which needs mount/unshare/pivot_root to build the sandbox (without @mount,
; bwrap fails with "Failed to make / slave"). The untrusted runner is then
; re-locked by the bwrap --seccomp denylist applied INSIDE the sandbox, so
; widening the trusted daemon's filter here does not widen the tenant code's.
SystemCallFilter=@system-service @mount
SystemCallErrorNumber=EPERM
; §41(a): delegate a cgroup-v2 subtree so the daemon can create per-unit child
; cgroups for the resource caps. Without it they degrade to the rlimit floor.
Delegate=yes
; Worker writes its state DB and (when libsecret is unavailable) the
; encrypted-file keystore under XDG_STATE_HOME / XDG_DATA_HOME — both under
; %h. ProtectHome=read-only would block those writes, so the standard XDG
; subdirectories are explicitly listed as read-write.
ReadWritePaths=%S/auspexai-worker %h/.local/state/auspexai-worker %h/.local/share/auspexai-worker

[Install]
WantedBy=default.target
