# packaging/systemd-deb/auspexai-worker.service
; AuspexAI worker — systemd user service unit, .deb install variant.
;
; This file is installed at /etc/systemd/user/auspexai-worker.service by
; the auspexai-worker .deb. Each volunteer user enables it with:
;
;   systemctl --user enable --now auspexai-worker.service
;
; The pip / uv install variant lives at packaging/systemd/auspexai-worker.service
; and points at %h/.local/bin/auspexai-worker; this one points at the
; .deb-installed venv at /opt/auspexai-worker/bin/auspexai-worker.
;
; Phase 1 hardening directives mirror the user-install variant; §5.17
; phasing applies equally to both.

[Unit]
Description=AuspexAI volunteer worker (Phase 1, .deb install)
Documentation=https://github.com/auspexai/worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/auspexai-worker/bin/auspexai-worker daemon
Restart=on-failure
RestartSec=10
; Low-priority scheduling — the worker runs in the volunteer's SPARE capacity and
; yields to their own foreground work, rather than competing at normal priority.
; Nice + a low cgroup CPU weight + the idle I/O class get out of the way under load.
Nice=19
CPUWeight=20
IOSchedulingClass=idle
; Hardening per §5.17 — applied at the daemon (trusted) tier. Phase 1
; permissive defaults; Phase 2 narrows further via config without
; altering the unit body.
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=strict
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
; §41(a): the daemon puts per-unit cgroup-v2 caps (memory.max/pids.max) around
; the runner, so the control-group tree must stay WRITABLE (=no) and delegated
; (Delegate=yes below). ProtectControlGroups=yes would remount it read-only.
ProtectControlGroups=no
RestrictRealtime=true
; @mount is REQUIRED on top of @system-service: the daemon spawns bubblewrap,
; which needs mount/unshare/pivot_root to build the sandbox (without @mount,
; bwrap fails with "Failed to make / slave"). The untrusted runner is then
; re-locked by the bwrap --seccomp denylist applied INSIDE the sandbox, so
; widening the trusted daemon's filter here does not widen the tenant code's.
SystemCallFilter=@system-service @mount
SystemCallErrorNumber=EPERM
; §41(a): delegate a cgroup-v2 subtree so the daemon can create per-unit child
; cgroups for the resource caps. Without it they degrade to the rlimit floor.
Delegate=yes
; Worker writes its state DB and the encrypted-file keystore fallback
; under XDG_STATE_HOME / XDG_DATA_HOME — both under %h. ProtectHome=
; read-only would block those writes, so the standard XDG subdirs are
; explicitly read-write. /opt/auspexai-worker/ is the venv install
; root and stays read-only.
ReadWritePaths=%S/auspexai-worker %h/.local/state/auspexai-worker %h/.local/share/auspexai-worker

[Install]
WantedBy=default.target
